# .config/systemd/user/volme3d.service
[Unit]
Description=Volme3D Webserver
After=network.target

[Service]
WorkingDirectory=/home/v3da
ExecStart=python3 /home/v3da/volme3d_server.py
Restart=always

[Install]
WantedBy=default.target
